# scripts/systemd/comandas-gateway.service
[Unit]
Description=Comandas API Gateway
After=network-online.target comandas-eureka.service comandas-auth.service comandas-producto.service comandas-pedido.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/comandas/spring-cloud-demo
EnvironmentFile=/etc/comandas/comandas.env
ExecStart=/usr/bin/java -Xms128m -Xmx384m -jar /opt/comandas/spring-cloud-demo/05-api-gateway/target/api-gateway-1.0.0.jar
Restart=on-failure
RestartSec=15

[Install]
WantedBy=multi-user.target
